# systemd/px-api-server.service
[Unit]
Description=PiCar-X REST API Server (:8420)
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/picar-x-hacking
Environment=HOME=/home/pi
ExecStartPre=/bin/sleep 2
ExecStart=/home/pi/picar-x-hacking/bin/px-api-server
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
